# server-package/esp32-server.service
[Unit]
Description=ESP32 Sensor Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/sketch_nov9a/server-package
ExecStart=/usr/bin/java -cp "esp32-server.jar:lib/*" com.esp32.server.EmbeddedServer 8085
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

